# deploy/launchd/dev.xiaoguai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd user agent for `xiaoguai serve` (T8.3).

  Rendered + installed by `xiaoguai service install` (macOS path): the
  BINARY / LOG_DIR placeholders (double-brace syntax) are substituted at
  install time with the resolved current executable and
  ~/Library/Logs/xiaoguai. Installed to
  ~/Library/LaunchAgents/dev.xiaoguai.plist — per-user, no root needed.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.xiaoguai</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{BINARY}}</string>
    <string>serve</string>
  </array>

  <!--
    #287: launchd agents do NOT inherit your shell's exports — API keys and
    XIAOGUAI_* settings set in ~/.zshrc are invisible to the service (so
    `xiaoguai doctor` can pass in your shell while the agent fails).
    Uncomment and fill in the dict below to pass environment to the service:

    <key>EnvironmentVariables</key>
    <dict>
      <key>OLLAMA_HOST</key>
      <string>http://127.0.0.1:11434</string>
      <key>XIAOGUAI_AUTH__USERNAME</key>
      <string>owner</string>
      <key>XIAOGUAI_AUTH__PASSWORD</key>
      <string>change-me</string>
    </dict>

    After editing the installed copy (~/Library/LaunchAgents/dev.xiaoguai.plist),
    reload it:  launchctl unload ~/Library/LaunchAgents/dev.xiaoguai.plist &&
                launchctl load -w ~/Library/LaunchAgents/dev.xiaoguai.plist
    NOTE: re-running `xiaoguai service install` overwrites the installed copy.
  -->

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>{{LOG_DIR}}/xiaoguai.out.log</string>

  <key>StandardErrorPath</key>
  <string>{{LOG_DIR}}/xiaoguai.err.log</string>
</dict>
</plist>
