# deploy/systemd/xiaoguai-core.service.d/10-low-port.conf
# Drop-in: bind to privileged ports (:80 / :443) without root.
#
# The base unit drops ALL capabilities (CapabilityBoundingSet=).
# Drop this file into the live drop-in directory to grant only the
# one capability xiaoguai-core needs for low-port binding, while
# keeping every other capability locked down.
#
# Install:
#   sudo cp 10-low-port.conf \
#       /etc/systemd/system/xiaoguai-core.service.d/10-low-port.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart xiaoguai-core
#
# Verify:
#   systemctl show xiaoguai-core | grep -E 'AmbientCapabilities|CapabilityBounding'
#
# To revert (return to all-caps-dropped default):
#   sudo rm /etc/systemd/system/xiaoguai-core.service.d/10-low-port.conf
#   sudo systemctl daemon-reload && sudo systemctl restart xiaoguai-core
#
# NOTE: You also need to change the bind address in config.yaml or
# $XIAOGUAI_SERVER__PORT to 80 or 443. The drop-in only grants the
# kernel permission to bind; it does not change what port the binary
# actually listens on.

[Service]
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
